# etc/services.d/mysqld/run
#!/bin/sh
mysqld --user=mysql --skip-name-resolve --skip-host-cache --skip-slave-start --debug-gdb